# gps/ostree/usr/lib/systemd/system/reunu-wakeup.service
[Unit]
Description=Reunu wakeup

[Service]
Type=oneshot
ExecStart=/usr/bin/reunu-enablemodem.sh

[Install]
WantedBy=suspend.target
